# deploy/short-telegram-bot-lite.service.example
[Unit]
Description=Telegram Short Signal Bot Lite
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/short-telegram-bot-lite
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/short-telegram-bot-lite/.venv/bin/python scripts/run_live.py
Restart=always
RestartSec=10
User=root

[Install]
WantedBy=multi-user.target
